# seedmelab/conf/drush.crontab
*/5 * * * * cd $DRUPAL_SITE_DIR && drush cron